# iotcaas.service
[Unit]
Description=IoT Laas Client Service

[Service]
ExecStart=/home/pi/start_at_boot/startiot_caas.sh

[Install]
WantedBy=multi-user.target
